# deploy/vetaris.service
[Unit]
Description=Vetaris E-Commerce Python Server
After=network.target

[Service]
# Kullanıcı adı ve grubu VPS'teki duruma göre değiştirebilirsiniz (örn: debian veya www-data)
User=debian
Group=debian

# Projenin bulunduğu dizin (VPS'te repo nereye çekildiyse)
WorkingDirectory=/var/www/vetaris.com

# Python sunucusunu başlatma komutu
ExecStart=/var/www/vetaris.com/venv/bin/python3 src/server.py

# Çökme durumunda otomatik yeniden başlat
Restart=always

[Install]
WantedBy=multi-user.target
